# bw_sys_initialize/scripts/rc.local
/sbin/ifconfig eth0 txqueuelen 1000
/sbin/ifconfig eth1 txqueuelen 1000